# systemd/daura-u.service.in
[Unit]
Description=Aura-U RGB Keyboard Daemon
After=multi-user.target

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/daura-u
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target